# infrastructure/managed-meeting/quadlet/chalk-redis.container.in
[Unit]
Description=Chalk ephemeral Redis acceleration
After=network-online.target
Wants=network-online.target

[Container]
Image=__REDIS_IMAGE_REFERENCE__
Pull=never
ContainerName=chalk-redis
Network=chalk-runtime.network
PodmanArgs=--network-alias=redis
Exec=redis-server --save '""' --appendonly no --maxmemory 256mb --maxmemory-policy allkeys-lru --protected-mode no --bind 0.0.0.0
User=999
Group=999
ReadOnly=true
ReadOnlyTmpfs=true
Tmpfs=/data:rw,nosuid,nodev,noexec,size=320m
DropCapability=all
NoNewPrivileges=true
PidsLimit=128
PodmanArgs=--memory=320m
HealthCmd=redis-cli ping
HealthInterval=10s
HealthRetries=3
HealthStartPeriod=5s
HealthTimeout=2s
HealthOnFailure=kill
Notify=healthy
LogDriver=journald
PodmanArgs=--log-opt=tag=chalk-redis
PodmanArgs=--stop-signal=SIGTERM
PodmanArgs=--stop-timeout=10

[Service]
Restart=on-failure
RestartSec=2s
StartLimitIntervalSec=60
StartLimitBurst=5
TimeoutStartSec=45
TimeoutStopSec=20

[Install]
WantedBy=chalk-runtime.target
