# infrastructure/managed-meeting/quadlet/chalk-cloudflared.container.in
[Unit]
Description=Chalk outbound-only Cloudflare Tunnel connector
After=network-online.target chalk-api.service chalk-sync.service
Requires=chalk-api.service chalk-sync.service
Wants=network-online.target

[Container]
Image=__CLOUDFLARED_IMAGE_REFERENCE__
Pull=never
ContainerName=chalk-cloudflared
Network=chalk-runtime.network
Secret=chalk-cloudflare-tunnel-token,type=mount,target=/run/secrets/tunnel-token,uid=65532,gid=65532,mode=0400
PodmanArgs=--entrypoint=/usr/local/bin/cloudflared
Exec=tunnel --no-autoupdate --loglevel info --metrics 0.0.0.0:2000 run --token-file /run/secrets/tunnel-token
PublishPort=127.0.0.1:12000:2000
User=65532
Group=65532
ReadOnly=true
ReadOnlyTmpfs=true
DropCapability=all
NoNewPrivileges=true
Ulimit=nofile=70000:70000
PidsLimit=256
PodmanArgs=--memory=384m
LogDriver=journald
PodmanArgs=--log-opt=tag=chalk-cloudflared
PodmanArgs=--stop-signal=SIGTERM
PodmanArgs=--stop-timeout=30

[Service]
Restart=on-failure
RestartSec=3s
StartLimitIntervalSec=120
StartLimitBurst=10
TimeoutStartSec=60
TimeoutStopSec=40

[Install]
WantedBy=chalk-runtime.target
