# infrastructure/managed-meeting/quadlet/chalk-api.container.in
[Unit]
Description=Chalk production API
After=network-online.target chalk-redis.service
Requires=chalk-redis.service
Wants=network-online.target

[Container]
Image=__API_IMAGE_REFERENCE__
Pull=never
ContainerName=chalk-api
Network=chalk-runtime.network
NetworkAlias=chalk-api
EnvironmentFile=/run/chalk/env/api.env
Secret=chalk-api-provider-server-cert,type=mount,target=/run/chalk/provider-bridge/server.crt,uid=65532,gid=65532,mode=0400
Secret=chalk-api-provider-server-key,type=mount,target=/run/chalk/provider-bridge/server.key,uid=65532,gid=65532,mode=0400
Secret=chalk-api-provider-client-ca,type=mount,target=/run/chalk/provider-bridge/client-ca.crt,uid=65532,gid=65532,mode=0400
PublishPort=127.0.0.1:18080:8080
User=65532
Group=65532
ReadOnly=true
ReadOnlyTmpfs=true
Tmpfs=/tmp:rw,nosuid,nodev,noexec,size=64m
DropCapability=all
NoNewPrivileges=true
PidsLimit=512
Memory=1g
HealthCmd=/usr/local/bin/chalk-healthcheck http://127.0.0.1:8080/readyz
HealthInterval=10s
HealthRetries=3
HealthStartPeriod=30s
HealthTimeout=3s
HealthOnFailure=kill
Notify=healthy
LogDriver=journald
LogOpt=tag=chalk-api
StopSignal=SIGTERM
StopTimeout=20

[Service]
Restart=on-failure
RestartSec=2s
StartLimitIntervalSec=60
StartLimitBurst=5
TimeoutStartSec=90
TimeoutStopSec=30

[Install]
WantedBy=chalk-runtime.target
