# infrastructure/managed-meeting/quadlet/chalk-sync.container.in
[Unit]
Description=Chalk production SyncEngine
After=network-online.target chalk-api.service
Requires=chalk-api.service
Wants=network-online.target

[Container]
Image=__SYNC_IMAGE_REFERENCE__
Pull=never
ContainerName=chalk-sync
Network=chalk-runtime.network
PodmanArgs=--network-alias=chalk-sync
EnvironmentFile=/run/chalk/env/sync.env
Secret=chalk-sync-provider-client-cert,type=mount,target=/run/chalk/provider-bridge/client.crt,uid=65532,gid=65532,mode=0400
Secret=chalk-sync-provider-client-key,type=mount,target=/run/chalk/provider-bridge/client.key,uid=65532,gid=65532,mode=0400
Secret=chalk-sync-provider-server-ca,type=mount,target=/run/chalk/provider-bridge/server-ca.crt,uid=65532,gid=65532,mode=0400
PublishPort=127.0.0.1:14100:4100
User=65532
Group=65532
ReadOnly=true
ReadOnlyTmpfs=true
Tmpfs=/tmp:rw,nosuid,nodev,noexec,size=128m
DropCapability=all
NoNewPrivileges=true
PidsLimit=2048
PodmanArgs=--memory=2g
HealthCmd=/usr/local/bin/chalk-healthcheck http://127.0.0.1:4100/readyz
HealthInterval=10s
HealthRetries=3
HealthStartPeriod=45s
HealthTimeout=3s
HealthOnFailure=kill
Notify=healthy
LogDriver=journald
PodmanArgs=--log-opt=tag=chalk-sync
PodmanArgs=--stop-signal=SIGTERM
PodmanArgs=--stop-timeout=15

[Service]
Restart=on-failure
RestartSec=2s
StartLimitIntervalSec=60
StartLimitBurst=5
TimeoutStartSec=120
TimeoutStopSec=25

[Install]
WantedBy=chalk-runtime.target
