# infrastructure/managed-meeting/systemd/chalk-runtime-watchdog.timer
[Unit]
Description=Run the Chalk host watchdog every 15 seconds

[Timer]
OnBootSec=45s
OnUnitActiveSec=15s
AccuracySec=1s
Unit=chalk-runtime-watchdog.service

[Install]
WantedBy=timers.target
